# systemd/rtltcp.service
[Unit]
Description=RTL TCP Service
After=network.target
Requires=rtltcp.socket

[Service]
Type=notify
ExecStart=/usr/local/bin/rtltcp
TimeoutStopSec=5
